# ops/systemd/hermes-root-executor.service.d/10-hermes-group-permissions.conf
[Service]
Group=hermes
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=restart
