# ops/systemd/hermes-root-executor.service.d/20-repository-commit.conf
[Service]
EnvironmentFile=/etc/hermes-root-executor/repository-commit.env
